# docker/supervisord-conf/init-db.conf
[program:init-db]
priority=2000
nodaemon=true
command=bash -c "sleep 15; /usr/bin/mysql < /www/database/MySQL/schema.sql"
