# docker/supervisord-conf/mysqld.conf
[program:mysql]
priority=1
startsecs=5
command=/usr/bin/pidproxy /var/run/mysqld/mysqld.pid /usr/bin/mysqld_safe
